# launchd/com.tradingview-watchlist-monitor.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tradingview-watchlist-monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/path/to/tradingview-watchlist-monitor/.venv/bin/python</string>
    <string>-u</string>
    <string>-m</string>
    <string>tv_watchlist_monitor.watcher</string>
    <string>daemon</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/path/to/tradingview-watchlist-monitor</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NOTIFIERS</key>
    <string>discord</string>
    <key>WATCHLIST_NAME</key>
    <string>TradingView Watchlist</string>
    <key>WATCHLIST_URL</key>
    <string>https://www.tradingview.com/watchlists/<watchlist-id>/</string>
    <key>STATE_FILE</key>
    <string>state/watchlist_<watchlist-id>.json</string>
    <key>MARKET_TIMEZONE</key>
    <string>America/New_York</string>
    <key>MARKET_OPEN</key>
    <string>09:30</string>
    <key>MARKET_CLOSE</key>
    <string>16:00</string>
    <key>CHECK_INTERVAL_MARKET_SECONDS</key>
    <string>900</string>
    <key>CHECK_INTERVAL_OFFHOURS_SECONDS</key>
    <string>3600</string>
    <key>SNAPSHOT_HOURS</key>
    <string>12</string>
    <key>SNAPSHOT_WINDOW_MINUTES</key>
    <string>45</string>
    <key>HEADLESS</key>
    <string>true</string>
    <key>DISCORD_WEBHOOK_URL</key>
    <string>PASTE_DISCORD_WEBHOOK_URL_HERE</string>
    <key>TELEGRAM_BOT_TOKEN</key>
    <string>PASTE_TELEGRAM_BOT_TOKEN_HERE</string>
    <key>TELEGRAM_CHAT_ID</key>
    <string>PASTE_TELEGRAM_CHAT_ID_HERE</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/path/to/tradingview-watchlist-monitor/state/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/path/to/tradingview-watchlist-monitor/state/launchd.err.log</string>
</dict>
</plist>
